# host/nftables/origin-firewall.service
# systemd unit that loads the Cloudflare-only origin firewall (nftables).
#
# Install (see docs/origin-firewall.md for the full, safe rollout):
#   sudo install -D -m 0644 host/nftables/origin-firewall.nft \
#       /etc/nftables.d/origin-firewall.nft
#   sudo install -D -m 0644 host/nftables/origin-firewall.service \
#       /etc/systemd/system/origin-firewall.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now origin-firewall.service
#
# The ruleset lives in its own `inet origin_fw` table, so this is independent of
# any distro nftables.service and only ever touches tcp/80,443.

[Unit]
Description=Cloudflare-only origin firewall (nftables)
# Load before k3s brings up klipper/kube-proxy so the web ports are guarded from
# the moment they can receive traffic. Independent of the main nftables.service.
Wants=network-pre.target
After=network-pre.target
Before=k3s.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/nft -f /etc/nftables.d/origin-firewall.nft
ExecReload=/usr/sbin/nft -f /etc/nftables.d/origin-firewall.nft
# On stop, remove ONLY our table — leave k3s/flannel/kube-proxy rules intact.
ExecStop=/usr/sbin/nft delete table inet origin_fw

[Install]
WantedBy=multi-user.target
